# tag-i3lock/config/systemd/user/screenshot.service
[Unit]
Description=Update screenshot every minute

[Service]
ExecStart=/home/stephen/.local/bin/screenshot.sh

[Install]
WantedBy=default.target
