# tag-i3lock/config/systemd/user/screenshot.timer
[Unit]
Description=Timer to run screenshot every minute

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
